# .docker/crond/etc/crontabs/www-data
* * * * * /usr/local/bin/php /data/www/bin/magento cron:run | grep -v "Ran jobs by schedule" >> /data/www/var/log/magento.cron.log
